# deploy/pillbox-telegram-bot.service
[Unit]
Description=Pillbox Telegram bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=asperelaiko
WorkingDirectory=/home/asperelaiko/pillbox
EnvironmentFile=/home/asperelaiko/pillbox/.env
ExecStart=/home/asperelaiko/pillbox/.venv/bin/python run_bot.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
